# logseq.conf
[program:logseq]
user=root
environment=HOME="/home/logseq",USER="logseq"
command=/app/run.sh
autorestart=true
